# scripts/browser-supervisor.conf
[program:chromium]
command=/opt/appliance/browser-launcher.sh
stdout_logfile=/tmp/chromium.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
autostart=true
autorestart=true
priority=40
stopasgroup=true
stopsignal=INT

[program:browser-api]
command=/opt/appliance/venv/bin/uvicorn app.server:app --host %(ENV_BROWSER_API_HOST)s --port %(ENV_BROWSER_API_PORT)s
stdout_logfile=/tmp/browser-api.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
directory=/opt/appliance
environment=PYTHONPATH="/opt/appliance"
autostart=true
autorestart=true
priority=50
stopasgroup=true
stopsignal=INT

[program:cdp-proxy]
command=/bin/sh -lc "exec socat TCP-LISTEN:%(ENV_BROWSER_CDP_PROXY_PORT)s,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:%(ENV_BROWSER_CDP_PORT)s"
stdout_logfile=/tmp/cdp-proxy.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
autostart=true
autorestart=true
priority=45
stopasgroup=true
stopsignal=INT
